# telegram-bot-api.service
[Unit]
Description=Telegram Bot API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app
Environment=TELEGRAM_API_ID=%i
Environment=TELEGRAM_API_HASH=%I
ExecStart=/usr/local/bin/telegram-bot-api --api-id %i --api-hash %I --local --http-port 8081 --log-level 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
